# files/etc/crontabs/root
41 4 */3 * * /opt/zapret/update-fake.sh >/dev/null 2>&1
0 4 1 * * /usr/bin/banned-guncelle

