# files/etc/uci-defaults/01-zapret
#!/bin/sh

/etc/init.d/zapret enable
/etc/init.d/zapret restart

exit 0