# files/etc/uci-defaults/01-zapret
#!/bin/sh
# /etc/uci-defaults/01-zapret

# Zapret init dosyasına izin ver
chmod +x /etc/init.d/zapret 2>/dev/null || true

# Zapret'i enable et (boot'ta çalışsın)
/etc/init.d/zapret enable 2>/dev/null || true

# Gerekli dizinleri oluştur
mkdir -p /opt/zapret/files/fake

# nfqws binary'e izin ver
chmod +x /opt/zapret/nfqws 2>/dev/null || true

exit 0